# scripts/launchd/com.baibot.pead.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  PEAD daily one-shot. Fires every weekday at 08:35 CDT (= 09:35 ET,
  ~5 min after market open). Auto-no-ops on weekends/holidays via the
  script's Alpaca-calendar check.

  Invokes scripts/launchd_wrappers/pead_run.sh which sources .env and
  execs the python runner with --live-submit. API keys NEVER appear in
  this plist file — they live in .env (gitignored).

  Install:
    cp scripts/launchd/com.baibot.pead.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.baibot.pead.plist

  Force a manual run NOW (still respects calendar — won't act on weekends):
    launchctl kickstart gui/$UID/com.baibot.pead

  Watch live:
    tail -f /Users/myu/code/baibot/results/pead/paper/launchd.out.log

  Remove:
    launchctl unload ~/Library/LaunchAgents/com.baibot.pead.plist
    rm ~/Library/LaunchAgents/com.baibot.pead.plist

  Strategy parameters are in the wrapper script. To change them
  (e.g. surprise threshold, hold days), edit
  scripts/launchd_wrappers/pead_run.sh and just kickstart again — no
  reload needed.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.baibot.pead</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/myu/code/baibot/scripts/launchd_wrappers/pead_run.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/myu/code/baibot</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>StartCalendarInterval</key>
  <array>
    <!-- Mon-Fri at 08:35 local (CDT/CST tracks the system clock automatically). -->
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>35</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>35</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>35</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>35</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>35</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/myu/code/baibot/results/pead/paper/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/myu/code/baibot/results/pead/paper/launchd.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
